# k8s/build/systemd/pelagos-install.service
# Triggered by pelagos-install.path on the guinea-pig node (ipc6). Installs the
# delivered binaries and restarts pelagos-cri. Oneshot — runs to completion each
# time the delivery drops a new .ready marker.
[Unit]
Description=Install staged pelagos binaries and restart pelagos-cri
After=pelagos-cri.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pelagos-install.sh
